# root/etc/periodic/daily/create_thumbnails
#!/usr/bin/with-contenv sh

LOGDIR=/config/logs/hb

su $USERNAME -c "cd /hb/ && bash -l -c 'python3 manage.py createthumbnails Pictures > $LOGDIR/createthumbs.log'"
