# root/etc/periodic/monthly/cleanup_thumbnails
#!/usr/bin/with-contenv sh

LOGDIR=/config/logs/hb

su $USERNAME -c "cd /hb/ && bash -l -c 'python3 manage.py thumbnail cleanup >> $LOGDIR/thumb_cleanup.log'"
